# meta-hulks/recipes-hulks/network-config/network-config/network-fallback-notification.timer
[Unit]
Description=Regularly start the notification

[Timer]
OnActiveSec=1m
